# src/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=10MB
pidfile=/tmp/supervisord.pid

; ── Uvicorn (FastAPI) ────────────────────────────────────────────────────────
[program:uvicorn]
command=uvicorn app.main:app --host 127.0.0.1 --port 8000 --workers 2
directory=/app/back-end
autostart=true
autorestart=true
startsecs=3
stdout_logfile=/var/log/supervisor/uvicorn.log
stderr_logfile=/var/log/supervisor/uvicorn.err.log
stdout_logfile_maxbytes=5MB
stderr_logfile_maxbytes=5MB

; ── Nginx ────────────────────────────────────────────────────────────────────
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx.err.log
stdout_logfile_maxbytes=5MB
stderr_logfile_maxbytes=5MB
